# etc/systemd/system/firefox_update.service.d/override.conf
[Service]
Environment="firefox_path=/opt/firefox/"
